# files/ipv6/lan1/bss/etc/uci-defaults/99-custom
uci -q batch << EOI
set network.lan.ipaddr='192.168.1.1'
set network.@device[0].igmp_snooping='1'
delete ucitrack.@cifs[-1]
add ucitrack cifs
set ucitrack.@cifs[-1].init=cifs
commit ucitrack
EOI
